# services/calypso.service
[unit]
Description=Calypso AI
[Service]
ExecStart=/usr/bin/calypso
[Install]
WantedBy=multi-user.target
